# backup/crontab
@daily /usr/bin/env python3 /app/backup_script.py | tee -a /app_data/backup.log
